# mynd.socket
[Unit]
Description=django mynd socket

[Socket]
ListenStream=/run/mynd.socket

[Install]
WantedBy=sockets.target
